# systemduser/.config/systemd/user/ludusavi.timer
[Unit]
Description="Ludusavi backup timer"

[Timer]
OnCalendar=weekly
Persistent=true
Unit=ludusavi.service

[Install]
WantedBy=timers.target

