# systemduser/.config/systemd/user/ludusavi.service
[Unit]
Description="Ludusavi backup"
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/bin/ludusavi backup --force
